# lib/systemd/system/networking.service.d/30_whonix-gw-firewall-fail-closed.conf
## Copyright (C) 2018 ENCRYPTED SUPPORT LP <adrelanos@riseup.net>
## See the file COPYING for copying conditions.

[Unit]
## Fail Closed Mechanism.
## When the Whonix firewall systemd service failed, do not bring up the
## network.
After=whonix-firewall.service
Requires=whonix-firewall.service
